# adws/deploy/adw-webhook.service
[Unit]
Description=ADW GitHub Webhook Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=adw-user
WorkingDirectory=/opt/adw
EnvironmentFile=/opt/adw/.env

# Use uv to run the webhook server
ExecStart=/opt/adw/.local/bin/uv run /opt/adw/adws/adw_triggers/trigger_webhook.py

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=0

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=adw-webhook

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Environment
Environment="PATH=/root/.cargo/bin:/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[Install]
WantedBy=multi-user.target
